# contrib/pulseaudio-audioclient.service
[Unit]
Description=Sound Service

# For audioclient, we always want to have pulseaudio started.

[Service]
# Note that notify will only work if --daemonize=no
Type=notify
#Environment=PULSE_CONFIG=/home/audioclient/.config/pulse/daemon.conf
ExecStart=/usr/bin/pulseaudio --daemonize=no
Restart=on-failure

[Install]
WantedBy=default.target
